# profile.d/20-aliases.sh
#!/usr/bin/env bash

# Easier navigation
alias ..="cd .."
alias ll="ls -hlF"
alias l.="ls -hlF -d .*"
alias la="ls -hla"

# Print public IP
alias pubip="dig +short myip.opendns.com @resolver1.opendns.com"

# Git
alias g="git"
alias gs="git status"
alias gc="git commit"
alias gd="git diff"
alias gdc="git diff --cached"
alias gl="git lp"

# Ruby
alias be="bundle exec"
