# profile.d/30-rbenv.sh
#!/usr/bin/env bash

if [ -d "$HOME/.rbenv" ]; then
    path_unshift "$HOME/.rbenv/bin"
    eval "$(rbenv init -)"
fi
